# .config/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /home/ubuntu/.config/djangobackend.ini

[program:nginx]
command=nginx

