# .config/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /home/ubuntu/django-sample-for-docker-compose/.config/djangobackend.ini

[program:nginx]
command=nginx

